# templates/systemd_teamcity.service.erb
[Unit]
Description=Teamcity master service

[Service]
Restart=on-failure
StartLimitInterval=20
StartLimitBurst=5
TimeoutStartSec=0
Environment="HOME=/var/lib/teamcity"
Environment="TEAMCITY_HOME=/var/lib/teamcity"
Environment="TEAMCITY_LOGS=/opt/teamcity/logs"
Environment="TEAMCITY_DATA_PATH=/var/lib/teamcity"

ExecStartPre=-/opt/teamcity/bin/teamcity-server stop 5 -force
ExecStart=/opt/teamcity/bin/<%= @start_command %>.sh run

ExecStop=-/opt/teamcity/bin/<%= @start_command %>.sh stop 30 -force

[Install]
WantedBy=multi-user.target
