# templates/systemd_teamcity.service.erb
[Unit]
Description=Teamcity master service

[Service]
Restart=on-failure
StartLimitInterval=20
StartLimitBurst=5
TimeoutStartSec=0
Environment="HOME=/var/lib/teamcity"
Environment="TEAMCITY_HOME=/var/lib/teamcity"
Environment="TEAMCITY_LOGS=/opt/teamcity/logs"
Environment="TEAMCITY_DATA_PATH=/var/lib/teamcity"

ExecStartPre=-/opt/teamcity/bin/<%= @kill_command %>
ExecStart=/opt/teamcity/bin/<%= @start_command %>

ExecStop=-/opt/teamcity/bin/<%= @stop_command %>

[Install]
WantedBy=multi-user.target
